# deploy/readiness-bot.service
[Unit]
Description=Readiness Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/readiness
Environment=TZ=America/Sao_Paulo
EnvironmentFile=/home/ubuntu/readiness/.env
ExecStart=/home/ubuntu/readiness/.venv/bin/python -m bot.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
